# service/data_collector.timer
[Unit]
Description=Collect Image and Sensor

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=false

[Install]
WantedBy=timers.target